# etc/update-motd.d/05-mahalia-banner
#!/bin/sh
cat /usr/share/mahalia-utils/motd
