# etc/update-motd.d/30-openmha-version
#!/bin/sh
mha --help | grep version
echo
